# services.d/cron/run
#!/command/with-contenv sh
exec supercronic /app/crontab
